# .config/systemd/user/backup-photos.service
[Service]
Type=oneshot
ExecStart=duply photos backup

[Unit]
Description=Backup Photo Library to B2 with Duply
