# .config/systemd/user/backup.timer
[Unit]
Description=Backup home to B2 with Duply (timer)

[Timer]
OnCalendar=daily
Unit=backup.service

[Install]
WantedBy=default.target
