# .config/systemd/user/backup-music.service
[Service]
Type=oneshot
ExecStart=duply music backup

[Unit]
Description=Backup Music Library to B2 with Duply
